# crontab
#STRART CRONTAB
* * * * * /usr/local/bin/python3 /app/main.py 
#END 